# video-converter.service
[Unit]
Description=Video Converter Daemon
After=network.target

[Service]
Type=simple
User=kb3kvq
Group=kb3kvq
WorkingDirectory=/home/kb3kvq/Documents/ClaudeProjects/VideoConverter
ExecStart=/usr/bin/python3 /home/kb3kvq/Documents/ClaudeProjects/VideoConverter/video_converter_daemon.py /home/kb3kvq/Documents/ClaudeProjects/VideoConverter/config.yaml
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
