# video-converter.service
[Unit]
Description=Video Converter Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/video-converter
ExecStart=/usr/bin/python3 /opt/video-converter/video_converter_daemon.py /opt/video-converter/config.yaml
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
